# deployment/minecraft-manager.service
[Unit]
Description=Minecraft Server Manager Web Application
After=network.target

[Service]
Type=simple
User=minecraft-manager
Group=minecraft-manager
WorkingDirectory=/opt/minecraft-manager
Environment=NODE_ENV=production
EnvironmentFile=/opt/minecraft-manager/.env
ExecStart=/usr/bin/node backend/dist/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target